# setup/ansible/templates/app.supervisor.conf.j2.ini
[program:api_server_{{ app_name }}]
user={{ user }}
directory={{ app_directory }}/back
command=bash -c "source {{ venv_directory }}/bin/activate && {{ venv_directory }}/bin/python{{ miniconda_python_version }} api_server.py"
autostart=true
autorestart=true
stdout_logfile={{ logs_directory }}/api_server.out.log
stderr_logfile={{ logs_directory }}/api_server.err.log

[program:rates_service_{{ app_name }}]
user={{ user }}
directory={{ app_directory }}/back
command=bash -c "source {{ venv_directory }}/bin/activate && {{ venv_directory }}/bin/python{{ miniconda_python_version }} rates_service.py"
autostart=true
autorestart=true
stdout_logfile={{ logs_directory }}/rates_service.out.log
stderr_logfile={{ logs_directory }}/rates_service.err.log

[program:tron_balances_updater_service_{{ app_name }}]
user={{ user }}
directory={{ app_directory }}/back
command=bash -c "source {{ venv_directory }}/bin/activate && {{ venv_directory }}/bin/python{{ miniconda_python_version }} tron_balances_updater_service.py"
autostart=true
autorestart=true
stdout_logfile={{ logs_directory }}/tron_balances_updater_service.out.log
stderr_logfile={{ logs_directory }}/tron_balances_updater_service.err.log

[program:db_backup_service_{{ app_name }}]
user={{ user }}
directory={{ app_directory }}/back
command=bash -c "source {{ venv_directory }}/bin/activate && {{ venv_directory }}/bin/python{{ miniconda_python_version }} db_backup_service.py"
autostart=true
autorestart=true
stdout_logfile={{ logs_directory }}/db_backup_service.out.log
stderr_logfile={{ logs_directory }}/db_backup_service.err.log

[group:{{ app_name }}]
programs=api_server_{{ app_name }}, rates_service_{{ app_name }}, tron_balances_updater_service_{{ app_name }}, db_backup_service_{{ app_name }}

